# services/MpiTouchBridge/run
#!/bin/sh
exec 2>&1
exec /usr/bin/python3 /data/victron-mpi5006-touch/touch-bridge.py
